# deploy/systemd/alphalens-thematic-daily.service
[Unit]
# AlphaLens thematic pipeline — runs the daily news → brief → JSON refresh
# end-to-end inside the alphalens-pipeline container. Triggered by the
# companion .timer. Web container picks up the new JSON via the shared
# bind-mount; nginx serves with no-cache so browsers see new briefs within
# minutes.
#
# Install as a *user* unit alongside the companion timer:
#
#   mkdir -p ~/.config/systemd/user
#   cp deploy/systemd/alphalens-thematic-daily.service ~/.config/systemd/user/
#   cp deploy/systemd/alphalens-thematic-daily.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now alphalens-thematic-daily.timer
#
# Inspect:           journalctl --user -u alphalens-thematic-daily.service -f
# Trigger manually:  systemctl --user start alphalens-thematic-daily.service
#
# Persistence (survives logout) — once, as the user:
#   sudo loginctl enable-linger "$USER"
Description=AlphaLens thematic pipeline daily run
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/AlphaLens

# Inherit jacoren UID/GID into the pipeline container so bind-mount writes
# land as the operator, not root. Compose reads UID/GID from the
# environment. Inject via `env` in ExecStart= because Environment= does not
# reliably expand %U/%G specifiers across systemd versions.
# --entrypoint /bin/bash overrides the image's
# ENTRYPOINT=["/app/.venv/bin/alphalens"] so the shell can execute the
# pipeline driver script. Without the override the script path would be
# passed to typer as a command name and fail with "No such command".
ExecStart=/usr/bin/env UID=%U GID=%G /usr/bin/docker compose \
    -f %h/AlphaLens/deploy/docker/docker-compose.yml \
    run --rm --entrypoint /bin/bash pipeline /app/deploy/docker/run_thematic_day.sh

# Keep stdout/stderr in journald.
StandardOutput=journal
StandardError=journal

# Pipeline runs Polygon news (~few MB JSON), Gemini API calls (network-bound),
# and pandas/parquet I/O. 4G is comfortable headroom on a small VPS; bump if
# the SP500-scale insider gate trips OOM.
MemoryMax=4G
TasksMax=512

# Oneshot doesn't auto-restart on success. On real failures (Polygon 5xx,
# Gemini quota, OOM) systemd surfaces non-zero exit in journalctl + status.
# The timer fires daily; a one-day miss is acceptable.

[Install]
WantedBy=default.target
