# deploy/systemd/alphalens-feedback-shadow-returns.service
[Unit]
# Nightly broker-free ladder + population-monitor backfill (market-behavior
# feedback). Replays each matured ledger decision's trade ladder over the last 14
# calendar days, then runs the population monitor over its own larger window —
# both deterministic price-path replays over Polygon minute bars (no broker).
# The legacy shadow-return / execution-quality metrics were removed with the
# broker chain (ADR 0012), and the Track-A user-action click ledger was removed
# later; the population-monitor replay (every brief candidate) is the primary
# market-behavior feedback that survives. The unit + command name are retained so
# the existing timer keeps working (a rename is a deferred follow-up). The Polygon
# minute-bar pull is rate-limited and the horizon must have closed (Polygon Basic
# serves only closed sessions).
#
# Install as a *user* unit alongside the companion timer:
#
#   mkdir -p ~/.config/systemd/user
#   cp deploy/systemd/alphalens-feedback-shadow-returns.service ~/.config/systemd/user/
#   cp deploy/systemd/alphalens-feedback-shadow-returns.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now alphalens-feedback-shadow-returns.timer
#
# Inspect:  journalctl --user -u alphalens-feedback-shadow-returns.service -f
# Trigger manually: systemctl --user start alphalens-feedback-shadow-returns.service
#
# Persistence (survives logout) — once, as the user:
#   sudo loginctl enable-linger "$USER"
Description=AlphaLens nightly broker-free ladder + population-monitor backfill
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/AlphaLens

# Shared secrets file: /etc/alphalens/env (root:root chmod 600). This job
# needs POLYGON_API_KEY (the minute-bar pricing leg). No leading `-`: fail
# loud on a missing/typoed path. A missing POLYGON_API_KEY does NOT hard-fail
# the run — every ticker fetch is skipped and the sweep reports "0 priced"
# (looks like a quiet night), so the loud-on-missing-file guard is the only
# protection against a silently mis-pointed env. See deploy/systemd/README.md
# "Environment file setup".
EnvironmentFile=/etc/alphalens/env

# Sweep mode: no --date arithmetic — the back window + per-date maturity guard
# pick the dates that need replaying. The ladder + population replays need no
# broker / account, so no --account flag.
ExecStart=%h/AlphaLens/.venv/bin/alphalens feedback backfill-shadow-returns

# Mirror the population-ladder parquet this recompute just (re)wrote into the
# edge Postgres cache via the maintenance ``rebuild-ladder-outcomes`` one-shot
# from the django-prod compose stack. This job is the SOLE writer of the
# parquet, so the mirror belongs HERE — it was moved off thematic-build, which
# re-synced unchanged data 5× per day and left the edge dashboard up to one
# slot (~2h) behind this recompute. ExecStartPost runs only after a successful
# ExecStart, so a failed recompute leaves the cache untouched. The spawned
# container runs under dockerd's cgroup, not this unit's, so it does not count
# against MemoryMax/TasksMax below. ``--profile maintenance`` opt-in isolates
# it from the always-on ``django`` service.
ExecStartPost=/usr/bin/docker compose \
    -f %h/AlphaLens/deploy/docker/django-prod/docker-compose.yaml \
    --profile maintenance run --rm rebuild-ladder-outcomes

# Prometheus textfile metrics (see alphalens-edgar-detect.service for the
# pattern). A clean run that prices 0 matured dates still counts as success
# for cron-health — last_success refreshes nightly so the AlphalensJobStale
# rule (48h) only fires when the sweep actually stops running.
ExecStopPost=%h/AlphaLens/deploy/systemd/bin/alphalens-emit-job-metrics feedback-shadow-returns

# Keep stdout/stderr in journald.
StandardOutput=journal
StandardError=journal

# The 14-day sweep x Polygon free-tier throttle (~5 req/min, 2 fetches per
# matured ticker) can run long after a VPS-downtime backlog — Persistent=true
# replay makes the first post-downtime run the largest. The recompute itself is
# idempotent, BUT the ExecStartPost Postgres mirror runs only after a SUCCESSFUL
# ExecStart — so a timeout-kill leaves fresh parquets on disk while the edge
# dashboard freezes at the last completed run's brief_date (it is NOT
# self-healing on its own). As the ongoing-position backlog grew (~100+ ongoing
# across ~30 brief dates, each deferring on the Polygon fetch budget), the
# steady-state run started crossing the old 45min ceiling and dying before the
# mirror, stranding the dashboard 2+ days behind. 90min gives headroom to reach
# ExecStartPost. Memory-light (minute-bar dicts), 1G is comfortable.
TimeoutStartSec=90min
MemoryMax=1G
TasksMax=32

[Install]
WantedBy=default.target
