# deploy/systemd/alphalens-paper-reconcile.service
[Unit]
# Paper-trade order-status sweep — fires every 30 min Mon-Fri during
# the XNYS regular session (13:30-20:00 UTC) plus a +30 / +60 min
# trailing window so final-minute fills get picked up. Companion to
# alphalens-paper-submit.service which queues the entry ladder pre-
# open. Both are PR-D of epic #295 (issue #298).
#
# Same two-layer gating as paper-submit: timer's ``OnCalendar=Mon..Fri``
# handles weekends; ``ExecCondition=`` below catches US public
# holidays falling on weekdays (~10/year).
#
# Reconcile does NOT take a ``--date`` arg — it sweeps all OPEN orders
# across all dates so a Friday-placed order that fills on Monday is
# still tracked. Adding ``--date $(date -u +%Y-%m-%d)`` would silently
# drop those cross-day orders.
#
# Install as a *user* unit alongside the companion timer:
#
#   mkdir -p ~/.config/systemd/user
#   cp deploy/systemd/alphalens-paper-reconcile.service ~/.config/systemd/user/
#   cp deploy/systemd/alphalens-paper-reconcile.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now alphalens-paper-reconcile.timer
#
# Inspect:  journalctl --user -u alphalens-paper-reconcile.service -f
# Trigger manually: systemctl --user start alphalens-paper-reconcile.service
Description=AlphaLens paper-trade order-status reconciler (Mon-Fri session)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/AlphaLens

# Same /etc/alphalens/env as paper-submit. Reconcile reads ALPACA_*
# only — no Telegram, no Polygon, no Gemini.
EnvironmentFile=/etc/alphalens/env

# ExecCondition gates on calendar reality — same rationale as paper-
# submit.service. Exit 0 proceeds, exit 1 skips silently.
ExecCondition=%h/AlphaLens/.venv/bin/alphalens paper is-trading-day

# NO --date flag — see top-of-file comment.
ExecStart=%h/AlphaLens/.venv/bin/alphalens paper reconcile

# Prometheus textfile metrics — same pattern as paper-submit.
ExecStopPost=%h/AlphaLens/deploy/systemd/bin/alphalens-emit-job-metrics paper-reconcile

StandardOutput=journal
StandardError=journal

# Reconcile makes O(N_open_orders) Alpaca REST calls; typically 5-30
# orders, ~5s wall time. 512M is comfortable; 1G matches the
# project-wide unit ceiling.
MemoryMax=1G
TasksMax=32

# A wedged Alpaca call blocking the next 30-min slot is the main risk.
# 5min = ~30× typical wall time; SIGTERM after that and the next slot
# takes over.
TimeoutStartSec=5min

[Install]
WantedBy=default.target
