# deploy/systemd/alphalens-paper-plan.timer
[Unit]
# Fires the companion alphalens-paper-plan.service Mon-Fri at 13:05
# UTC — after the 12:30 UTC thematic-build wrote that morning's
# (D-1) brief parquet (normal finish ~12:45-12:58 UTC) and 20 min
# before paper-submit fires at 13:25 UTC.
#
# **Why 13:05**: leaves margin on both sides — enough slack past a
# slow build, enough lead before submit reads the PLANNED rows. Do
# not push past ~13:15 or the margin before submit shrinks.
#
# **Why Mon..Fri not daily**: paper-submit only runs Mon-Fri, so
# planning on a non-submit day would write rows nothing consumes. The
# companion ExecCondition catches US public holidays that slip
# through Mon..Fri.
#
# **Why NOT Persistent=true**: a weekend/holiday catch-up fire would
# plan against a stale brief. Same rule as paper-submit / paper-
# reconcile — no catch-up of missed slots; the next trading day's
# fire takes over.
Description=Trigger AlphaLens paper-trade plan Mon-Fri pre-submit
After=network-online.target

[Timer]
OnCalendar=Mon..Fri *-*-* 13:05:00 UTC
Unit=alphalens-paper-plan.service

[Install]
WantedBy=timers.target
