# deploy/systemd/alphalens-paper-plan.service
[Unit]
# Paper-trade daily planning — the first link in the chain. Fires once
# per US trading day at 13:05 UTC, reads that morning's brief parquet,
# and writes PLANNED ladder rows into the SQLite ledger. The companion
# alphalens-paper-submit.service (13:25 UTC) then pushes those PLANNED
# rows into Alpaca paper; alphalens-paper-reconcile.service sweeps the
# fills. Without this unit the ledger has no PLANNED rows and submit
# silently submits nothing. PR-D follow-up of epic #295 (issue #298).
#
# **Date contract (load-bearing).** `alphalens thematic brief` defaults
# to (today-1) and run_thematic_day.sh runs it with no --date, so the
# 12:30 UTC build on day D writes (D-1).parquet — the freshest complete
# brief. `plan --date` is required and reads briefs_dir/<date>.parquet;
# submit re-reads the ledger keyed on the same brief_date. So plan AND
# submit MUST pass --date yesterday: a mismatch means plan reads a
# missing file (FileNotFoundError) or submit finds zero PLANNED rows.
#
# **Two-layer gating** (same as submit/reconcile):
#   1. The companion timer's `OnCalendar=Mon..Fri ...` filters weekends.
#   2. `ExecCondition=` runs `alphalens paper is-trading-day` to catch
#      US public holidays falling on weekdays. NOTE: unlike submit/
#      reconcile, the plan subcommand has NO internal market-closed
#      guard, so this ExecCondition is the ONLY holiday filter for plan.
#
# Install as a *user* unit alongside the companion timer:
#
#   cp deploy/systemd/alphalens-paper-plan.service ~/.config/systemd/user/
#   cp deploy/systemd/alphalens-paper-plan.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now alphalens-paper-plan.timer
#
# Inspect:  journalctl --user -u alphalens-paper-plan.service -f
# Trigger manually: systemctl --user start alphalens-paper-plan.service
Description=AlphaLens paper-trade daily planning (Mon-Fri pre-submit)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/AlphaLens

# /etc/alphalens/env carries ALPACA_API_KEY + ALPACA_API_SECRET (paper-
# main) and ALPACA_TEST_* for the dev sandbox. No leading `-`: a
# missing key MUST fail loud rather than silently plan against the
# wrong account.
EnvironmentFile=/etc/alphalens/env

# Holiday gate — the only one for plan (no internal market-closed
# guard in the subcommand). Exit 0 proceeds, exit 1 skips silently.
ExecCondition=%h/AlphaLens/.venv/bin/alphalens paper is-trading-day

# Plan YESTERDAY UTC's brief (see date contract above). `--use-test-
# account` routes to the Alpaca TEST paper account during the testing
# phase, matching the submit/reconcile units — plan tags PLANNED rows
# with this account and submit re-reads by (brief_date, account), so
# all three MUST agree. NO `--force`: it deletes existing plans +
# shadow_log for the brief_date, and a re-fire would wipe the rows
# submit is about to read. Wrapped in `/bin/sh -c '...'` for the
# `$(date)` substitution; `%%` is systemd's literal-percent escape.
ExecStart=/bin/sh -c '%h/AlphaLens/.venv/bin/alphalens paper plan --use-test-account --date $(date -u -d yesterday +%%Y-%%m-%%d)'

# Prometheus textfile metrics — same pattern as the other active units.
# Runs unconditionally so a failed plan still updates last_exit_code.
ExecStopPost=%h/AlphaLens/deploy/systemd/bin/alphalens-emit-job-metrics paper-plan

StandardOutput=journal
StandardError=journal

# Plan does a handful of Alpaca REST reads (equity) + local parquet/
# ledger work. 1G matches the project-wide unit ceiling.
MemoryMax=1G
TasksMax=32

# Type=oneshot defaults TimeoutStartSec to infinity. A plan that wedges
# on Alpaca's API would block tomorrow's fire. 5min = ample headroom.
TimeoutStartSec=5min

[Install]
WantedBy=default.target
