# deploy/systemd/alphalens-thematic-daily.service
[Unit]
# AlphaLens thematic pipeline — runs the daily news → brief refresh
# end-to-end via the pipeline image (built standalone from
# deploy/docker/Dockerfile.pipeline). The Django stack picks up the
# new parquet output via a one-shot `rebuild-cache` invocation in
# ExecStartPost.
#
# Install as a *user* unit alongside the companion timer:
#
#   mkdir -p ~/.config/systemd/user
#   cp deploy/systemd/alphalens-thematic-daily.service ~/.config/systemd/user/
#   cp deploy/systemd/alphalens-thematic-daily.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now alphalens-thematic-daily.timer
#
# Inspect:           journalctl --user -u alphalens-thematic-daily.service -f
# Trigger manually:  systemctl --user start alphalens-thematic-daily.service
#
# Persistence (survives logout) — once, as the user:
#   sudo loginctl enable-linger "$USER"
Description=AlphaLens thematic pipeline daily run
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/AlphaLens

# Inherit operator UID/GID into the pipeline container so bind-mount writes
# land as the user, not root.
ExecStart=/usr/bin/env UID=%U GID=%G /usr/bin/docker run --rm \
    --user %U:%G \
    -e HOME=/app/home \
    -v %h/.alphalens:/app/home/.alphalens \
    --entrypoint /bin/bash \
    alphalens-pipeline:latest \
    /app/deploy/docker/run_thematic_day.sh

# Sync the parquet output the pipeline just wrote into the Django Postgres
# cache via the maintenance one-shot. ``--profile maintenance`` opt-in
# isolates this from the always-on ``django`` service. ExecStartPost fires
# only on ExecStart success; pipeline failures leave the cache untouched.
ExecStartPost=/usr/bin/docker compose \
    -f %h/AlphaLens/deploy/docker/django-prod/docker-compose.yaml \
    --profile maintenance run --rm rebuild-cache

# Keep stdout/stderr in journald.
StandardOutput=journal
StandardError=journal

# Pipeline runs Polygon news (~few MB JSON), Gemini API calls (network-bound),
# and pandas/parquet I/O. 4G is comfortable headroom on a small VPS; bump if
# the SP500-scale insider gate trips OOM.
MemoryMax=4G
TasksMax=512

# Oneshot doesn't auto-restart on success. On real failures (Polygon 5xx,
# Gemini quota, OOM) systemd surfaces non-zero exit in journalctl + status.
# The timer fires daily; a one-day miss is acceptable.

[Install]
WantedBy=default.target
