# deploy/systemd/alphalens-feedback-shadow-returns.service
[Unit]
# Nightly broker-free population-monitor backfill (market-behavior feedback).
# Runs the population monitor (every brief candidate) over its own multi-week
# window — deterministic price-path replay over Polygon minute bars (no broker).
# The legacy shadow-return / execution-quality metrics were removed with the
# broker chain (ADR 0012), and the Track-A user-action click ledger was removed
# later (#465/#467) — the per-decision ladder replay that read it is gone too,
# and the old ~/.alphalens/feedback.db file (if present on a host) is an orphan
# that nothing opens at runtime. The population-monitor replay is now the sole
# market-behavior feedback signal (parquet + Postgres edge_ladderoutcome). The unit + command name are retained so
# the existing timer keeps working (a rename is a deferred follow-up). The Polygon
# minute-bar pull is rate-limited and the horizon must have closed (Polygon Basic
# serves only closed sessions).
#
# Install as a *user* unit alongside the companion timer:
#
#   mkdir -p ~/.config/systemd/user
#   cp deploy/systemd/alphalens-feedback-shadow-returns.service ~/.config/systemd/user/
#   cp deploy/systemd/alphalens-feedback-shadow-returns.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now alphalens-feedback-shadow-returns.timer
#
# Inspect:  journalctl --user -u alphalens-feedback-shadow-returns.service -f
# Trigger manually: systemctl --user start alphalens-feedback-shadow-returns.service
#
# Persistence (survives logout) — once, as the user:
#   sudo loginctl enable-linger "$USER"
Description=AlphaLens nightly broker-free ladder + population-monitor backfill
After=network-online.target docker.service
Wants=network-online.target
# The edge mirror (rebuild-ladder-outcomes) has been decoupled into its own
# unit (alphalens-edge-mirror.service). OnSuccess= fires it immediately after a
# successful nightly recompute so /edge refreshes right away. On a failure or
# 90-min timeout-kill, the mirror's own hourly timer (the "self-heal backstop"
# in alphalens-edge-mirror.service) re-syncs whatever parquets exist (including
# any partial output from the killed run) within <=1h — a latency-only trade-off
# — so we deliberately do NOT also set OnFailure= here. Pointing BOTH OnSuccess=
# and OnFailure= at the
# same target made systemd register two identical trigger sources on edge-mirror
# and log a per-run "multiple trigger source candidates for exit status
# propagation (X, X), skipping" warning; dropping OnFailure= removes the warning
# without losing failure-path coverage. See deploy/systemd/README.md
# "Edge mirror (decoupled)" for deployment atomic requirements.
OnSuccess=alphalens-edge-mirror.service

[Service]
Type=oneshot
WorkingDirectory=%h/AlphaLens

# Shared secrets file: /etc/alphalens/env (root:root chmod 600). This job
# needs POLYGON_API_KEY (the minute-bar pricing leg). No leading `-`: fail
# loud on a missing/typoed path. A missing POLYGON_API_KEY does NOT hard-fail
# the run — every ticker fetch is skipped and the sweep reports "0 priced"
# (looks like a quiet night), so the loud-on-missing-file guard is the only
# protection against a silently mis-pointed env. See deploy/systemd/README.md
# "Environment file setup".
EnvironmentFile=/etc/alphalens/env

# Sweep mode: no --date arithmetic — the back window + per-date maturity guard
# pick the dates that need replaying. The ladder + population replays need no
# broker / account, so no --account flag.
ExecStart=%h/AlphaLens/.venv/bin/alphalens feedback backfill-shadow-returns

# Prometheus textfile metrics (see alphalens-edgar-detect.service for the
# pattern). A clean run that prices 0 matured dates still counts as success
# for cron-health — last_success refreshes nightly so the AlphalensJobStale
# rule (48h) only fires when the sweep actually stops running.
ExecStopPost=%h/AlphaLens/deploy/systemd/bin/alphalens-emit-job-metrics feedback-shadow-returns

# Keep stdout/stderr in journald.
StandardOutput=journal
StandardError=journal

# The 14-day sweep x Polygon free-tier throttle (~5 req/min, 2 fetches per
# matured ticker) can run long after a VPS-downtime backlog — Persistent=true
# replay makes the first post-downtime run the largest. As the ongoing-position
# backlog grew (~100+ ongoing across ~30 brief dates), the steady-state run
# started crossing the old 45min ceiling. 90min gives the nightly sweep headroom
# to complete. The edge Postgres mirror (alphalens-edge-mirror.service) is a
# separate unit fired via OnSuccess= after a clean run; a timeout-kill is caught
# by the mirror's own hourly self-heal timer within <=1h (see the OnSuccess=
# note in [Unit] above). Memory-light (minute-bar dicts), 1G is
# comfortable.
TimeoutStartSec=90min
MemoryMax=1G
TasksMax=32

[Install]
WantedBy=default.target
