# deploy/systemd/av-earnings-backfill.service
[Unit]
# Alpha Vantage EARNINGS daily backfill — oneshot triggered by the companion
# .timer. Each fire consumes up to the AV free-tier 25-call/day quota, exits
# clean, and waits for tomorrow.
#
# Install as a *user* unit alongside the companion timer:
#
#   mkdir -p ~/.config/systemd/user
#   cp deploy/systemd/av-earnings-backfill.service ~/.config/systemd/user/
#   cp deploy/systemd/av-earnings-backfill.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now av-earnings-backfill.timer
#
# Inspect:  journalctl --user -u av-earnings-backfill.service -f
# Trigger manually: systemctl --user start av-earnings-backfill.service
#
# Persistence (survives logout) — once, as the user:
#   sudo loginctl enable-linger "$USER"
Description=AlphaLens Alpha Vantage EARNINGS daily backfill
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/AlphaLens

# `.env` must contain ALPHA_VANTAGE_API_KEY=...; chmod 600 recommended.
# Systemd reads each KEY=VALUE line from this file before ExecStart.
EnvironmentFile=%h/AlphaLens/.env

ExecStart=%h/AlphaLens/.venv/bin/python -m scripts.av_earnings_daily_backfill \
    --universe sp500_union

# Keep stdout/stderr in journald (no /var/log file needed).
StandardOutput=journal
StandardError=journal

# AV client is memory-light (JSON dict per ticker, max 25/day). 1G is
# comfortable headroom; bump if a future SP1500-scale run trips OOM.
MemoryMax=1G
TasksMax=32

# Oneshot does not auto-restart on success-with-clean-exit (AVRateLimitError
# path), which is the expected steady-state behaviour. Real failures (network
# storm, malformed AV response) should surface in journalctl rather than be
# silently retried — the timer fires daily so a one-day skip is non-fatal.

[Install]
WantedBy=default.target
